# packages/device-scanner-daemon/systemd-units/device-scanner.service
[Unit]
Description=MFL Device Scanner Daemon
RefuseManualStart=true
DefaultDependencies=false

[Service]
Restart=always
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /sbin/device-scanner-daemon
