# IML.ScannerProxyDaemon/systemd-units/scanner-proxy.path
[Unit]
Description=IML Scanner Proxy Path

[Path]
DirectoryNotEmpty=/var/lib/chroma/settings

[Install]
WantedBy=default.target
