# packages/block-device-listener/udev-rules/99-iml-device-scanner.rules
ACTION=="add|remove", SUBSYSTEM=="block", ENV{IML_SIZE}="$attr{size}", RUN+="/root/dist/bundle.js"
