# packages/block-device-listener/udev-rules/99-iml-device-scanner.rules
ACTION=="add|remove", SUBSYSTEM=="block", ENV{MFL_SIZE}="$attr{size}", RUN+="/lib/udev/block-device-listener"
